# deploy/systemd/rt-forwarder.service
# deploy/systemd/rt-forwarder.service
#
# Systemd unit for the rt-forwarder service.
#
# Installation:
#   sudo cp deploy/systemd/rt-forwarder.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable rt-forwarder
#   sudo systemctl start rt-forwarder
#
# Configuration:
#   /etc/rusty-timer/forwarder.toml   — forwarder config (TOML, required)
#   /var/lib/rusty-timer/             — journal and token storage
#
# See docs/runbooks/forwarder-operations.md for full operational guidance.

[Unit]
Description=Remote Timing Forwarder (rt-forwarder)
Documentation=https://github.com/iwismer/rusty-timer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Service identity — run as a dedicated non-root user.
User=rt-forwarder
Group=rt-forwarder

# Binary and config location.
ExecStart=/usr/local/bin/rt-forwarder

# Data and config directories.
WorkingDirectory=/var/lib/rusty-timer

# Environment — forwarder reads TOML config, not env vars.
# RUST_LOG can be set here for log verbosity override (e.g. debug, trace).
Environment=RUST_LOG=info

# Restart policy: always restart on failure, 5-second delay between attempts.
# This handles transient network failures and system reboots.
Restart=on-failure
RestartSec=5s

# Limit restart loops: maximum 5 restarts in 60 seconds before giving up.
StartLimitInterval=60s
StartLimitBurst=5

# Logging: journal output captured by journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rt-forwarder

# Security hardening.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/rusty-timer
ReadWritePaths=/etc/rusty-timer

# Give the forwarder time to flush the SQLite WAL before SIGKILL.
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
